# system/usr/local/etc/rc.d/cloudflared
#!/bin/sh

# PROVIDE: cloudflared
# REQUIRE: NETWORKING SERVERS
# KEYWORD: shutdown

. /etc/rc.subr

name="cloudflared"
rcvar="cloudflared_enable"

procname="/usr/local/bin/cloudflared"
pidfile="/var/run/cloudflared.pid"
logfile="/var/log/cloudflared.log"
token_file="/usr/local/etc/cloudflared/token"

: ${cloudflared_enable:="NO"}
# For token-based tunnels, "tunnel run" is the expected form.
: ${cloudflared_mode:="tunnel run"}
# Handy default; override in rc.conf if you want updates.
: ${cloudflared_flags:="--no-autoupdate"}

load_rc_config $name

# Build cloudflared args (append token if present)
cf_args="${cloudflared_mode} ${cloudflared_flags}"
if [ -f "${token_file}" ]; then
    cf_args="${cf_args} --token $(cat "${token_file}")"
fi

command="/usr/sbin/daemon"
# -r = respawn if child exits
command_args="-r -o ${logfile} -p ${pidfile} ${procname} ${cf_args}"

start_precmd="${name}_prestart"
cloudflared_prestart()
{
    mkdir -p /var/run
    mkdir -p "$(dirname "${logfile}")"

    case "${cloudflared_mode}" in
        *"tunnel run"*)
            if [ ! -f "${token_file}" ]; then
                echo "${name}: missing ${token_file} (required for 'tunnel run')."
                return 1
            fi
            ;;
    esac
}

run_rc_command "$1"